# rpi_server/rpi-health.service
[Unit]
Description=RPI Health Node — publishes CPU/temp/memory metrics to foxglove-bridge
After=foxglove-bridge.service dronepi-hotspot.service network.target
Wants=foxglove-bridge.service

[Service]
Type=simple
User=dronepi
WorkingDirectory=/home/dronepi/unitree_lidar_project/rpi_server
Environment=HOME=/home/dronepi
ExecStart=/home/dronepi/unitree_lidar_project/rpi_server/start_rpi_health.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rpi-health

[Install]
WantedBy=multi-user.target
