# unitree_drone_mapper/flight/drone-watchdog.service
[Unit]
Description=DronePi Flight Stack Watchdog (Point-LIO + bag recorder)
# Must start after MAVROS so the /mavros/state topic exists to subscribe to.
# Also depends on the SSD mount so bag recording has a valid output path.
After=mavros.service mnt-ssd.mount dronepi-hotspot.service
Wants=mavros.service mnt-ssd.mount dronepi-hotspot.service

[Service]
Type=simple
User=dronepi
Environment="HOME=/home/dronepi"

# Source ROS 2 before running the watchdog so rclpy and mavros_msgs
# are importable. The watchdog itself launches Point-LIO and the bag
# recorder as subprocesses with their own source calls.
ExecStart=/bin/bash -c "\
    source /opt/ros/jazzy/setup.bash && \
    source /home/dronepi/unitree_lidar_project/RPI5/ros2_ws/install/setup.bash && \
    python3 /home/dronepi/unitree_lidar_project/unitree_drone_mapper/flight/drone_watchdog.py"

# Restart the watchdog if it crashes -- but not if it exits cleanly (e.g.
# systemctl stop). on-failure avoids restart loops on unrecoverable errors.
Restart=on-failure
RestartSec=10

# Allow time for the watchdog to stop Point-LIO and close the bag cleanly
# before systemd sends SIGKILL. Bag close takes up to 5s (GRACEFUL_KILL_S).
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
